# cds/systemd/cds-master.service
[Unit]
# CDS Master systemd unit — P4 Part 18 hardening edition.
#
# Installs CDS as a system service so the kernel restarts it after
# any crash (OOM, segfault, self-update failures). Pairs with
# /api/self-update's pre-restart validation: if an operator triggers
# self-update with broken code, the validation step aborts BEFORE
# killing this process. If the process is somehow killed anyway,
# systemd Restart=always brings it back.
#
# Install (recommended): use the helper command which auto-fills
# the paths for your install location:
#
#   cd cds && ./exec_cds.sh install-systemd
#
# Install (manual):
#   # 1. Edit WorkingDirectory/Environment/ExecStart paths below
#   sudo cp cds/systemd/cds-master.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now cds-master
#
# Inspect:
#   systemctl status cds-master
#   journalctl -u cds-master -f
#
# See doc/design.cds-resilience.md Phase 2 + the P4 Part 18 hardening
# commit for why ExecStartPre exists.

Description=CDS Master (Cloud Development Suite)
Documentation=file:///home/user/prd_agent/doc/design.cds-resilience.md
After=docker.service network-online.target
Wants=network-online.target
Requires=docker.service

# Rate limit: max 5 restarts in 60 seconds. These directives MUST be
# in [Unit], not [Service], per current systemd. Putting them in
# [Service] silently no-ops and emits a warning in `systemd-analyze
# verify`. Caught during P4 Part 18 self-verify pass.
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple

# IMPORTANT: edit these paths to match your install location.
# `exec_cds.sh install-systemd` does this for you automatically.
WorkingDirectory=/opt/prd_agent/cds
# systemd does not inherit the interactive shell PATH. Keep node's bin
# directory on PATH so pnpm/npx shebangs (`/usr/bin/env node`) can resolve.
# `exec_cds.sh install-systemd` rewrites this line to include nvm/asdf paths
# when Node is installed outside /usr/bin.
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=NODE_ENV=production
Environment=CDS_REPO_ROOT=/opt/prd_agent
# Optionally pin the AI access key here (prefer /etc/cds/env for secrets)
EnvironmentFile=-/etc/cds/env

# ──────────────────────────────────────────────────────────────
# P4 Part 18 hardening — ExecStartPre self-heal(2026-05-06 减负版)。
#
# 用户实测 self-update 215s,定位发现 ~50s 来自 systemd restart 时
# **再跑一次** `npx tsc`(emit 到 dist),与 self-update 流程已经做的
# tsc 完全冗余。self-force-sync 还会清 dist/.tsbuildinfo,导致 systemd
# 这次 tsc 走全量冷启动 30-50s。
#
# 现行策略:
#   - pnpm install --frozen-lockfile:保留(快路径 ~5s,自动拉新依赖)
#   - npx tsc:**删除**(self-update 路径自己负责 dist,systemd 不再编译)
#
# Fail-safe(self-update 失败 / 误删 dist):node 启动会报
# "Cannot find module dist/index.js",systemd 标 failed,operator 在
# journalctl 看到错误后用 /api/self-force-sync 重新构建。比之前"systemd
# 帮你救场"的代价低 30-50s/次,且大多数路径不会失败(self-update 自带
# 预检 + 原子 swap)。
# ──────────────────────────────────────────────────────────────
ExecStartPre=/usr/bin/pnpm install --frozen-lockfile

ExecStart=/usr/bin/node dist/index.js

# Restart policy — the whole point of this unit file.
# (The matching StartLimitIntervalSec / StartLimitBurst rate-limit
# is in [Unit] above per modern systemd convention.)
Restart=always
RestartSec=3s

# 启动超时：默认 90s 不够 —— self-update / self-force-sync 会清空 dist/，
# 下次 systemd ExecStartPre 全量 `pnpm install + npx tsc` 冷启动可能要
# 2-3 分钟。给 5 分钟兜底，避免"删了 dist 重启就启动超时"的事故。
TimeoutStartSec=300

# Standard output goes to the journal
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cds-master

# Resource accounting (requires cgroup v2)
#
# 2026-05-06 Bugbot be8dc4f5 / 2c22c57a 反馈 + 用户授权:
# 之前 MemoryMax=512M 是 master 自身 + 它 spawn 的所有 child(tsc / vite / pnpm)
# 共享的 cgroup 上限。我们已经移除了 NODE_OPTIONS=--max-old-space-size,V8 默认
# 会按主机 RAM 自适应,但 V8 看不到 cgroup 限制 → 跑大型 tsc/vite 时会被
# SIGKILL(uncatchable,无 graceful 错误)。
#
# 用户明确策略:"每个容器都不限制内存,尽情释放" + "不要刻意放大"。
# MemoryMax 解除(infinity = 无限制),让 master + child 跟随主机 RAM。
# CPUQuota 保留作为防爆 host 的兜底。
MemoryMax=infinity
CPUQuota=100%

# Security hardening — soft, compatible with docker CLI usage
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/opt/prd_agent /var/run/docker.sock
ProtectHome=false

[Install]
WantedBy=multi-user.target
