# cds/systemd/cds-master.service
[Unit]
# CDS Master systemd unit — P4 Part 18 hardening edition.
#
# Installs CDS as a system service so the kernel restarts it after
# any crash (OOM, segfault, self-update failures). Pairs with
# /api/self-update's pre-restart validation: if an operator triggers
# self-update with broken code, the validation step aborts BEFORE
# killing this process. If the process is somehow killed anyway,
# systemd Restart=always brings it back.
#
# Install (recommended): use the helper command which auto-fills
# the paths for your install location:
#
#   cd cds && ./exec_cds.sh install-systemd
#
# Install (manual):
#   # 1. Edit WorkingDirectory/Environment/ExecStart paths below
#   sudo cp cds/systemd/cds-master.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now cds-master
#
# Inspect:
#   systemctl status cds-master
#   journalctl -u cds-master -f
#
# See doc/design.cds-resilience.md Phase 2 + the P4 Part 18 hardening
# commit for why ExecStartPre exists.

Description=CDS Master (Cloud Development Suite)
Documentation=file:///home/user/prd_agent/doc/design.cds-resilience.md
After=docker.service network-online.target
Wants=network-online.target
Requires=docker.service

# Rate limit: max 5 restarts in 60 seconds. These directives MUST be
# in [Unit], not [Service], per current systemd. Putting them in
# [Service] silently no-ops and emits a warning in `systemd-analyze
# verify`. Caught during P4 Part 18 self-verify pass.
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple

# IMPORTANT: edit these paths to match your install location.
# `exec_cds.sh install-systemd` does this for you automatically.
WorkingDirectory=/opt/prd_agent/cds
# systemd does not inherit the interactive shell PATH. Keep node's bin
# directory on PATH so pnpm/npx shebangs (`/usr/bin/env node`) can resolve.
# `exec_cds.sh install-systemd` rewrites this line to include nvm/asdf paths
# when Node is installed outside /usr/bin.
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=NODE_ENV=production
Environment=CDS_REPO_ROOT=/opt/prd_agent
# Optionally pin the AI access key here (prefer /etc/cds/env for secrets)
EnvironmentFile=-/etc/cds/env

# ──────────────────────────────────────────────────────────────
# 2026-05-06 重构(用户反馈"为什么每次都要手动 sudo 重装 unit"):
# 把"启动 master 进程要做什么"全部委托给 exec_cds.sh 的 master-run 子命令。
# 本 unit 文件只声明**框架**(在哪儿跑、Restart 策略、资源限制),所有
# 易变的启动细节(pnpm 标志、node 启动参数、ExecStartPre 步骤)随
# self-update 自动更新到 exec_cds.sh,operator 不用每次都 sudo cp。
#
# unit 文件改动**只**在以下情况发生(应该极少):
#   - 改 systemd 行为本身(Restart 策略、MemoryMax、TimeoutStartSec)
#   - 改 WorkingDirectory(install 位置变了)
# 其它启动细节去 cds/exec_cds.sh 的 `master-run` case 改即可。
# ──────────────────────────────────────────────────────────────
# install-systemd 时会把这条路径替换为实际安装位置 + 实际用户家目录的 PATH
ExecStart=/opt/prd_agent/cds/exec_cds.sh master-run

# Restart policy — the whole point of this unit file.
# (The matching StartLimitIntervalSec / StartLimitBurst rate-limit
# is in [Unit] above per modern systemd convention.)
Restart=always
RestartSec=3s

# 启动超时：默认 90s 不够 —— self-update / self-force-sync 会清空 dist/，
# 下次 systemd ExecStartPre 全量 `pnpm install + npx tsc` 冷启动可能要
# 2-3 分钟。给 5 分钟兜底，避免"删了 dist 重启就启动超时"的事故。
TimeoutStartSec=300

# Standard output goes to the journal
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cds-master

# Resource accounting (requires cgroup v2)
#
# 2026-05-06 Bugbot be8dc4f5 / 2c22c57a 反馈 + 用户授权:
# 之前 MemoryMax=512M 是 master 自身 + 它 spawn 的所有 child(tsc / vite / pnpm)
# 共享的 cgroup 上限。我们已经移除了 NODE_OPTIONS=--max-old-space-size,V8 默认
# 会按主机 RAM 自适应,但 V8 看不到 cgroup 限制 → 跑大型 tsc/vite 时会被
# SIGKILL(uncatchable,无 graceful 错误)。
#
# 用户明确策略:"每个容器都不限制内存,尽情释放" + "不要刻意放大"。
# MemoryMax 解除(infinity = 无限制),让 master + child 跟随主机 RAM。
# CPUQuota 保留作为防爆 host 的兜底。
MemoryMax=infinity
CPUQuota=100%

# Security hardening — soft, compatible with docker CLI usage
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/opt/prd_agent /var/run/docker.sock
ProtectHome=false

[Install]
WantedBy=multi-user.target
